# etc/sv/certificate-issuer/run
#!/bin/sh
exec chpst -uroot /usr/bin/certificate-issuer 2>&1
