# etc/sv/ntpd/run
#!/bin/sh

exec chpst -uroot /usr/sbin/ntpd -n 2>&1
