# etc/sv/ntpd/run
#!/bin/sh

exec chpst -uroot /usr/sbin/ntpd 2>&1
