# etc/sv/agetty/run
#!/bin/sh

exec chpst -uroot /sbin/agetty --keep-baud 115200,38400,9600 -a root ttyS1 linux 2>&1
