# etc/sv/ic-gateway/run
#!/bin/sh
exec 2>&1

. /mnt/ic-gateway/ic-gateway.env

exec chpst -uroot /usr/bin/ic-gateway
