# etc/sv/ic-gateway/run
#!/bin/sh
set -a

if [ ! -d /mnt/ic-gateway ]; then
    exit
fi

. /mnt/ic-gateway/ic-gateway.env
ulimit -n 262144
exec chpst -u root /usr/bin/ic-gateway 2>&1
