# etc/sv/ic-gateway/run
#!/bin/sh
exec chpst -uroot /usr/bin/ic-gateway 2>&1
